# deploy/youtube-stream.service
[Unit]
Description=Hummingbird Cam YouTube Live Stream
After=hummingbird-cam.service network-online.target
Requires=hummingbird-cam.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/opt/hummingbird-cam/stream-youtube.sh
Restart=on-failure
RestartSec=30

# Run as motion user
User=motion
Group=motion

# Stream key loaded from env file
EnvironmentFile=/etc/hummingbird-cam.env

[Install]
# Timer-activated — no WantedBy so it won't auto-start on boot
