# deploy/youtube-stream-start.timer
[Unit]
Description=Start Hummingbird YouTube stream at 7 AM daily

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
